# com.moatasim.enhanced-document-organization.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.moatasim.enhanced-document-organization</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/moatasimfarooque/Downloads/Programming/CascadeProjects/Drive_sync/drive_sync.sh</string>
        <string>all</string>
    </array>
    
    <!-- Run every hour (3600 seconds) -->
    <key>StartInterval</key>
    <integer>3600</integer>
    
    <!-- Run at system startup/login -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Don't keep the job alive if it exits -->
    <key>KeepAlive</key>
    <false/>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/moatasimfarooque</string>
    </dict>
    
    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/moatasimfarooque/Downloads/Programming/CascadeProjects/Drive_sync</string>
    
    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/moatasimfarooque/Downloads/Programming/CascadeProjects/Drive_sync/logs/automation.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/moatasimfarooque/Downloads/Programming/CascadeProjects/Drive_sync/logs/automation_error.log</string>
    
    <!-- Only run when user is logged in -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
    
</dict>
</plist>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
    
    <!-- Process type -->
    <key>ProcessType</key>
    <string>Background</string>
    
    <!-- Resource limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>1024</integer>
    </dict>
    
    <!-- Nice value (lower priority) -->
    <key>Nice</key>
    <integer>10</integer>
    
    <!-- Throttle interval between restarts -->
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>