# com.moatasim.drive_sync.file_watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.moatasim.drive_sync.file_watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/moatasimfarooque/Downloads/Programming/CascadeProjects/Drive_sync/mcp/mcp_client.js</string>
        <string>watch_files</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/moatasimfarooque/Downloads/Programming/CascadeProjects/Drive_sync/logs/file_watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/moatasimfarooque/Downloads/Programming/CascadeProjects/Drive_sync/logs/file_watcher_error.log</string>
</dict>
</plist>